# package/vibeops.service
[Unit]
Description=vibeopsOS First Boot Setup
After=network.target network-online.target
Wants=network.target network-online.target
ConditionPathExists=!/etc/vibeops-initialized

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/vibeops/vibeops-init.sh

[Install]
WantedBy=multi-user.target
